# init-scripts/abrt-dump-location.path
[Unit]
Description=ABRT Dump Location Watch
After=systemd-tmpfiles-setup.service abrtd.socket
Before=paths.target

[Path]
Unit=abrtd.service
PathChanged=/var/tmp/abrt
MakeDirectory=yes
