# init-scripts/abrt-pstoreoops.service
[Unit]
Description=Collect UEFI-saved oopses for ABRT
After=abrtd.socket
Requisite=abrtd.socket
ConditionDirectoryNotEmpty=/sys/fs/pstore

[Service]
Type=oneshot
ExecStart=/usr/sbin/abrt-harvest-pstoreoops
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
